# rules/40-redhat-hotunplug.rules
# ACPI0004 container offline for Huawei Kunlun
# do not edit this file, it will be overwritten on update

SUBSYSTEM=="container", ACTION=="change", DEVPATH=="*/ACPI0004:??", \
RUN+="/bin/sh -c ' \
if [ $(cat /sys/$env{DEVPATH}/online) -eq 1 ]; then \
        find -L /sys/$env{DEVPATH}/firmware_node/*/physical_node* -maxdepth 1 -name online | \
        while read line; do \
                if [ $(cat $line) -eq 1 ]; then \
                        /bin/echo 0 > $line; \
                fi \
        done; \
        /bin/echo 0 > /sys/$env{DEVPATH}/online; \
fi'"